# clamav/clamav.crontab
0 9 0 0 0 /usr/bin/bash /home/ec2-user/clamav/run-clamav.sh