# clamav/clamav.crontab
0 9 * * * /usr/bin/bash /home/astun/clamav/run-clamav.sh
